# services/latency_monitor.service
[Unit]
Description=RTSP camera network discovery
After=network.target

[Service]
Type=simple
WorkingDirectory=/app/rtsp_streamer
ExecStart=/app/rtsp_streamer/scripts/latency_monitor.py
Restart=always
RestartSec=1
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
